# luci-app-kickass/root/etc/init.d/kickass
#!/bin/sh /etc/rc.common
START=50

start()
{
	exec bash /usr/share/kickass/kickass 1 &
}

stop()
{
	echo "1" >> /tmp/kickass/stop
	sleep 1
	rm -rf /tmp/kickass/stop
}

restart()
{
	stop
	start
}
